# .config/systemd/user/ssh-agent.service
[Unit]
Description=SSH key agent

[Service]
Type=simple
ExecStart=/usr/bin/ssh-agent -D -a /tmp/ssh-agent.socket

[Install]
WantedBy=default.target
